# files/shared/etc/profile.d/horizon-aliases.sh
#!/usr/bin/env bash

# All the aliases

alias ls="eza"
alias cat="bat -p"
alias fastfetch="/usr/bin/fastfetch --config /usr/share/horizon/configs/fastfetch/config.jsonc"
